# heaven/tools/com.ym.claude-transcript-archiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ym.claude-transcript-archiver</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ym/Developer/llm/heaven/tools/archive_old_transcripts.py</string>
        <string>--apply</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/ym/.claude/transcript-archiver.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ym/.claude/transcript-archiver.err.log</string>
</dict>
</plist>
